# etc/profile.d/cs50.sh
# aliases
alias cp="cp -i"
alias ll="ls -l --color=auto"
alias mv="mv -i"
alias rm="rm -i"

# environment variables
export EDITOR=nano
